# linux/systemd/smart-monitor.service
[Unit]
Description=Smart Monitor Service
Documentation=https://github.com/smart-monitor/smart-monitor
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/smart_monitor
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
User=root
Group=root

# Environment variables
Environment=RUST_LOG=info
Environment=GST_DEBUG=3

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
